# etc/profile.d/golang.sh
export GOROOT=$DOTSOFTWARE_ROOT_DIR/opt/golang
export PATH=$PATH:$GOROOT/bin
